# supabase-keepalive/systemd/supabase-keepalive.service
[Unit]
Description=Supabase Keep-Alive - Prevents free-tier project pausing
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/rasppi-utils/supabase-keepalive
ExecStart=/opt/rasppi-utils/.venv/bin/python /opt/rasppi-utils/supabase-keepalive/scripts/keepalive.py
EnvironmentFile=/etc/rasppi-utils/supabase-keepalive/.env

[Install]
WantedBy=multi-user.target
