# social-poster/systemd/social-poster-tunnel.service
[Unit]
Description=Cloudflare Tunnel for social-poster (public image URLs so Instagram can fetch uploads)
After=network-online.target social-poster.service
Wants=network-online.target

[Service]
Type=simple
# Runs the named tunnel defined in /etc/cloudflared/config.yml. That config and
# its credentials are created by `social-poster/setup-cloudflare.sh` (one-time,
# needs a browser authorization). Until then this restart-loops harmlessly;
# Bluesky posting works without it.
ExecStart=/usr/bin/cloudflared tunnel --no-autoupdate --config /etc/cloudflared/config.yml run
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
